# server/delchef.upstart
# del-chef upstart script
description "Starts and stops the delchef server"
version "1.0"
author "Max Dietz"

start on networking

respawn

exec /usr/bin/delchef
